# alpine3.7-apache2.4-php5.6/root/etc/sv/memcached/run
#!/bin/sh
#exec chpst -u nobody /usr/bin/memcached -m 64
exec /usr/bin/memcached  -u root -m 64
